# hippo-sync-worker.service
[Unit]
Description=Hippo Sync Worker Service
After=network.target postgresql.service

[Service]
Type=simple
User=sab
WorkingDirectory=/home/sab/hippo
ExecStart=/home/sab/hippo/.venv/bin/python -m hippo sync-worker
Environment=HIPPO_LOG_FILE=/home/sab/hippo/logs/hippo.log
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
